# dot_config/systemd/user/gpg-forwarded-sockets.path
[Unit]
Description=Monitor gpg forward socket 

[Path]
PathExists=/tmp/%u.gpg.forward.socket
Unit=gpg-forwarded-sockets.service

[Install]
WantedBy=default.target
