# dot_config/systemd/user/plasma-plasmashell.service.d/gpg-agent-reload.conf
[Service]
ExecStartPost=/usr/bin/systemctl --user restart gpg-agent.service
