# dot_config/exact_systemd/exact_user/exact_gpg-agent-ssh.socket.d/env.conf
[Socket]
ExecStopPost=-/usr/bin/systemctl --user unser-environment SSH_AUTH_SOCK
ExecStartPost=-/usr/bin/systemctl --user set-environment SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh
