# dot_config/exact_systemd/user/podman.socket.d/env.conf
[Socket]
ExecStartPost=-/usr/bin/systemctl --user set-environment DOCKER_HOST=unix://%t/podman/podman.sock
