# dot_config/exact_systemd/user/gpg-agent-ssh.socket.d/env.conf
[Socket]
ExecStartPost=-/usr/bin/systemctl --user set-environment DOCKER_HOST=%t/gnupg/S.gpg-agent.ssh
