# dot_config/exact_systemd/user/gpg-agent-ssh.socket.d/env.conf
[Socket]
ExecStartPost=-/home/joke2/.local/bin/debug /usr/bin/systemctl --user set-environment TEST_JOKE_PASS=234234
